# systemd/myapp.service
[Unit]
Description=My App
After=network.target


[Service]
ExecStart=/usr/local/bin/myapp.sh

[Install]
WantedBy=multi-user.target
