# systemd/myapp.service
[Unit]
Description=My Awesome App
After=network.target


[Service]
ExecStart=/root/myapp.sh

[Install]
WantedBy=multi-user.target
